# cron/crontab
0 0 * * * mysqldump -uvagrant -pvagrant mysql | gzip > /vagrant/db/dumps/$(date +\%Y\%m\%d\%H\%M\%S).sql.gz
